# debian/cocaine-tornado-proxy.upstart
description "cocaine-tornado-proxy"

start on filesystem
stop on runlevel [016]
limit nofile 65535 65535

pre-start script
    if [ -f /etc/default/$UPSTART_JOB ]; then
        . /etc/default/$UPSTART_JOB
    fi
end script

script
    if [ -f /etc/default/$UPSTART_JOB ]; then
        . /etc/default/$UPSTART_JOB
    fi

   exec sudo -u cocaine /usr/bin/cocaine-tornado-proxy --config=$CONFIG_PATH
end script
